# iso/provisioning/cortex-first-boot.service
[Unit]
Description=Cortex Linux First-Boot Provisioning
Documentation=https://github.com/cortexlinux/cortex-distro
After=network-pre.target systemd-resolved.service
Before=network.target getty@tty1.service display-manager.service
Wants=network-pre.target
ConditionPathExists=!/opt/cortex/provisioning/.first-boot-complete

[Service]
Type=oneshot
ExecStart=/opt/cortex/provisioning/first-boot.sh
RemainAfterExit=yes
TimeoutSec=600

# Security hardening for the service itself
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false
PrivateTmp=true

# Logging
StandardOutput=journal+console
StandardError=journal+console
SyslogIdentifier=cortex-first-boot

[Install]
WantedBy=multi-user.target
