# packages/cortex-branding/etc/update-motd.d/00-cortex-banner
#!/bin/bash
#
# Cortex Linux MOTD Banner
# Displays branded welcome message on terminal login
#

# Colors
PURPLE='\033[0;35m'
CYAN='\033[0;36m'
WHITE='\033[1;37m'
GRAY='\033[0;90m'
RESET='\033[0m'

# ASCII Art Logo
cat << 'EOF'

                 ██████╗ ██████╗ ██████╗ ████████╗███████╗██╗  ██╗
                ██╔════╝██╔═══██╗██╔══██╗╚══██╔══╝██╔════╝╚██╗██╔╝
                ██║     ██║   ██║██████╔╝   ██║   █████╗   ╚███╔╝
                ██║     ██║   ██║██╔══██╗   ██║   ██╔══╝   ██╔██╗
                ╚██████╗╚██████╔╝██║  ██║   ██║   ███████╗██╔╝ ██╗
                 ╚═════╝ ╚═════╝ ╚═╝  ╚═╝   ╚═╝   ╚══════╝╚═╝  ╚═╝

EOF

echo -e "${GRAY}                              L I N U X${RESET}"
echo ""
