# packages/cx-branding/etc/update-motd.d/10-cx-sysinfo
#!/bin/bash
#
# CX Linux System Information
# Displays system stats on terminal login
#

# Colors
PURPLE='\033[0;35m'
CYAN='\033[0;36m'
WHITE='\033[1;37m'
GRAY='\033[0;90m'
GREEN='\033[0;32m'
YELLOW='\033[0;33m'
RED='\033[0;31m'
RESET='\033[0m'

# Gather system information
HOSTNAME=$(hostname)
DISTRO=$(grep PRETTY_NAME /etc/os-release 2>/dev/null | cut -d'"' -f2 || echo "CX Linux")
KERNEL=$(uname -r)
UPTIME=$(uptime -p 2>/dev/null | sed 's/up //' || uptime | awk -F'( |,|:)+' '{print $6"h "$7"m"}')
LOAD=$(cat /proc/loadavg | awk '{print $1" "$2" "$3}')

# CPU info
CPU_MODEL=$(grep "model name" /proc/cpuinfo 2>/dev/null | head -1 | cut -d':' -f2 | sed 's/^ //' | cut -c1-40)
CPU_CORES=$(nproc 2>/dev/null || echo "?")

# Memory info
if command -v free &>/dev/null; then
    MEM_TOTAL=$(free -h | awk '/^Mem:/ {print $2}')
    MEM_USED=$(free -h | awk '/^Mem:/ {print $3}')
    MEM_PERCENT=$(free | awk '/^Mem:/ {printf "%.0f", $3/$2*100}')
else
    MEM_TOTAL="N/A"
    MEM_USED="N/A"
    MEM_PERCENT="0"
fi

# Disk info
if command -v df &>/dev/null; then
    DISK_TOTAL=$(df -h / | awk 'NR==2 {print $2}')
    DISK_USED=$(df -h / | awk 'NR==2 {print $3}')
    DISK_PERCENT=$(df / | awk 'NR==2 {print $5}' | tr -d '%')
else
    DISK_TOTAL="N/A"
    DISK_USED="N/A"
    DISK_PERCENT="0"
fi

# Network info
IP_ADDR=$(hostname -I 2>/dev/null | awk '{print $1}' || echo "N/A")

# Users logged in
USERS=$(who 2>/dev/null | wc -l || echo "0")

# Color based on usage
mem_color() {
    if [ "$1" -lt 50 ]; then echo -e "${GREEN}";
    elif [ "$1" -lt 80 ]; then echo -e "${YELLOW}";
    else echo -e "${RED}"; fi
}

disk_color() {
    if [ "$1" -lt 70 ]; then echo -e "${GREEN}";
    elif [ "$1" -lt 90 ]; then echo -e "${YELLOW}";
    else echo -e "${RED}"; fi
}

# Display
echo -e "  ${GRAY}┌──────────────────────────────────────────────────────────────────┐${RESET}"
echo -e "  ${GRAY}│${RESET}  ${PURPLE}System${RESET}     ${WHITE}${HOSTNAME}${RESET}"
echo -e "  ${GRAY}│${RESET}  ${PURPLE}OS${RESET}         ${DISTRO}"
echo -e "  ${GRAY}│${RESET}  ${PURPLE}Kernel${RESET}     ${KERNEL}"
echo -e "  ${GRAY}│${RESET}  ${PURPLE}Uptime${RESET}     ${UPTIME}"
echo -e "  ${GRAY}│${RESET}  ${PURPLE}Load${RESET}       ${LOAD}"
echo -e "  ${GRAY}├──────────────────────────────────────────────────────────────────┤${RESET}"
echo -e "  ${GRAY}│${RESET}  ${CYAN}CPU${RESET}        ${CPU_MODEL:-Unknown} (${CPU_CORES} cores)"
echo -e "  ${GRAY}│${RESET}  ${CYAN}Memory${RESET}     $(mem_color $MEM_PERCENT)${MEM_USED}${RESET} / ${MEM_TOTAL} (${MEM_PERCENT}%)"
echo -e "  ${GRAY}│${RESET}  ${CYAN}Disk${RESET}       $(disk_color $DISK_PERCENT)${DISK_USED}${RESET} / ${DISK_TOTAL} (${DISK_PERCENT}%)"
echo -e "  ${GRAY}│${RESET}  ${CYAN}IP${RESET}         ${IP_ADDR}"
echo -e "  ${GRAY}│${RESET}  ${CYAN}Users${RESET}      ${USERS} logged in"
echo -e "  ${GRAY}└──────────────────────────────────────────────────────────────────┘${RESET}"
echo ""
