# packages/cx-branding/etc/update-motd.d/99-cx-footer
#!/bin/bash
#
# CX Linux MOTD Footer
# Shows helpful links and tips
#

# Colors
PURPLE='\033[0;35m'
CYAN='\033[0;36m'
GRAY='\033[0;90m'
RESET='\033[0m'

# Random tips
TIPS=(
    "Use 'cx' command for AI-powered shell assistance"
    "Run 'cx-update' to check for CX component updates"
    "Enable dark mode: gsettings set org.gnome.desktop.interface color-scheme prefer-dark"
    "Check security status: sudo lynis audit system"
    "View system logs: journalctl -xe"
    "Monitor resources: htop or btop"
)

# Pick random tip
TIP="${TIPS[$RANDOM % ${#TIPS[@]}]}"

echo -e "  ${GRAY}─────────────────────────────────────────────────────────────────────${RESET}"
echo -e "  ${PURPLE}Tip:${RESET} ${TIP}"
echo -e "  ${GRAY}─────────────────────────────────────────────────────────────────────${RESET}"
echo -e "  ${CYAN}Docs:${RESET} https://docs.cxlinux.com  ${CYAN}Support:${RESET} https://cxlinux.com/support"
echo ""
